# a2a/agent_zero/kitt-agent.service
[Unit]
Description=KITT Agent Zero (Sovereign Administrator)
After=network.target docker.service

[Service]
User=<operator>
Group=<operator>
WorkingDirectory=/home/<operator>/kitt_gateway/a2a/agent_zero
ExecStart=/home/<operator>/kitt_gateway/a2a/agent_zero/venv/bin/python /home/<operator>/kitt_gateway/a2a/agent_zero/agent.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target