# a2a/agent_zero/kitt-agent-http.service
[Unit]
Description=KITT Agent Zero HTTP Service
After=network.target docker.service

[Service]
User=doubl
WorkingDirectory=/home/doubl/kitt_gateway/a2a/agent_zero
EnvironmentFile=/home/doubl/kitt_gateway/a2a/agent_zero/.env
ExecStart=/home/doubl/kitt_gateway/a2a/agent_zero/venv/bin/python agent.py serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
