# hub/kitt-hub.service
[Unit]
Description=KITT Hub — Sovereign AI Gateway
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=doubl
WorkingDirectory=/home/doubl/kitt_gateway/hub
ExecStart=/home/doubl/kitt_gateway/venv/bin/python main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
